# winnerbot.service
[Unit]
Description=Post to yourewinner.com

[Service]
Type=simple
ExecStart=/home/steven/src/markov/start.sh
