# winnerbot.timer
[Unit]
Description=Post to yourewinner.com daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
